# code/python.service
[Unit]

[Service]
ExecStart=/bin/sh -c "nohup python3 /tmp/code/hello.py"